# .config/systemd/user/numlock-osd.service

[Unit]
Description=Num Lock OSD Notification Service

[Service]
Type=simple
ExecStart=%h/.local/bin/numlock-osd.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target