# .config/systemd/user/volume-osd.service
[Unit]
Description=Volume OSD Notification Service

[Service]
Type=simple
ExecStart=%h/.local/bin/volume-osd.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target